# scripts/com.tommyroar.tailscale-auth.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.tommyroar.tailscale-auth</string>
    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/node</string>
        <string>/Users/tommydoerr/dev/tommyroar.github.io/scripts/auth_server.js</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>WorkingDirectory</key>
    <string>/Users/tommydoerr/dev/tommyroar.github.io</string>
    <key>StandardOutPath</key>
    <string>/Users/tommydoerr/dev/tommyroar.github.io/logs/auth_server.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/tommydoerr/dev/tommyroar.github.io/logs/auth_server.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
</dict>
</plist>
